# debian-latest/misc/crontab
@reboot sed -i 's/127.0.1.1 ip-/127.0.1.1 NAME ip-/g' /etc/hosts
@reboot hostnamectl set-hostname NAME
